# Dev/Spark-Cassandra.vsphere-16.04/provision/roles/hadoop/templates/hdfs-namenode.service.j2
[Unit]
Description=Hadoop Distributed File System NameNode

[Service]
Type=forking
ExecStart={{ hadoop_home }}/sbin/hadoop-daemon.sh --script hdfs start namenode
User={{ cluster_user }}

[Install]
WantedBy=multi-user.target