# Test/test.dse.vsphere/provision/roles/admin/templates/__upstart.conf.j2
description "jupyter"

start on runlevel [2345]
stop on shutdown

respawn
respawn limit 10 15

setuid {{ cluster_user }}
setgid {{ cluster_user }}

exec ipython notebook --no-browser
